# autotrade-backend/deploy/systemd/autotrade-celery-worker.service
[Unit]
Description=AutoTrade Pro — Celery worker (mark-to-market, journal sync, scans)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/cis/windows/auto-trade-pro/autotrade-backend
Environment="PYTHONPATH=/home/cis/windows/auto-trade-pro/autotrade-backend"
ExecStart=/home/cis/windows/auto-trade-pro/autotrade-backend/.venv/bin/celery -A tasks.celery_app worker --loglevel=warning --concurrency=2
Restart=always
RestartSec=5
# Restart even if it crashes repeatedly
StartLimitIntervalSec=0
StandardOutput=append:/tmp/celery_worker.log
StandardError=append:/tmp/celery_worker.log

[Install]
WantedBy=default.target
